# s6-rc.d/rpcbind/run
#!/bin/bash

echo "启动 RPC 服务"
exec /sbin/rpcbind -f
